# dataset_platform_linux/dataset_platform.service
# systemd 服务示例
# 复制到 /etc/systemd/system/dataset_platform.service 后，
# 将 YOUR_USERNAME 改为实际 Linux 用户名（若解压到 Documents，路径如下）。
# 若使用 Conda，请将 ExecStart 改为 conda 环境中的 streamlit 绝对路径。
#
# 启用：sudo systemctl daemon-reload && sudo systemctl enable dataset_platform && sudo systemctl start dataset_platform

[Unit]
Description=数据集标注平台 (Streamlit)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/Documents/dataset_platform_linux
Environment="PATH=/home/YOUR_USERNAME/Documents/dataset_platform_linux/venv/bin"
Environment="DATA_ROOT=/home/YOUR_USERNAME/Documents/dataset_platform_linux/dataset"
ExecStart=/home/YOUR_USERNAME/Documents/dataset_platform_linux/venv/bin/streamlit run dataset_viewer/app.py --server.address 0.0.0.0 --server.port 8501
Restart=on-failure

[Install]
WantedBy=multi-user.target
